# scripts/systemd/crypto-watchdog.service
# Copie de reference de l'unite VPS (watchdog, alerte seule) - capturee le 2026-07-16
# pour le runbook de restauration (docs/runbook-restauration-vps.md).
[Unit]
Description=Crypto AI Terminal — Watchdog (surveillance moteur, alerte seule)
After=network-online.target crypto-advisor.service
Wants=network-online.target

[Service]
Type=simple
User=mathieuhasard111
WorkingDirectory=/home/mathieuhasard111/crypto_ai_terminal
EnvironmentFile=-/home/mathieuhasard111/crypto_ai_terminal/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8
Environment=TZ=UTC
ExecStart=/home/mathieuhasard111/crypto_ai_terminal/.venv/bin/python3 watchdog_vps.py
Restart=on-failure
RestartSec=15s
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto_watchdog
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
